# templates/default/init/consul.service.erb
[Unit]
Description="HashiCorp Consul - A service mesh solution"
Documentation=https://www.consul.io/
Requires=network-online.target
After=network-online.target
ConditionFileNotEmpty=<%= node['consul']['conf_dir'] %>/consul.hcl

[Service]
EnvironmentFile=<%= node['consul']['conf_dir'] %>/systemd_env_vars
User=<%= node['consul']['user'] %>
Group=<%= node['consul']['group'] %>
ExecStart=<%= node['consul']['bin_dir'] %>/consul agent -bind <%= @bind_address %> -config-dir=<%= node['consul']['conf_dir'] %>
ExecReload=<%= node['consul']['bin_dir'] %>/consul reload
KillMode=process
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target